# rootfs/etc/crontabs/caddy
0 8 * * * /usr/bin/php -f /app/artisan ninja:send-invoices > /dev/null 2>&1
0 8 * * * /usr/bin/php -f /app/artisan ninja:send-reminders > /dev/null 2>&1
0 8 * * * /usr/bin/php -f /app/artisan ninja:send-renewals > /dev/null 2>&1
